# deploy/systemd/upstream-probe.service
[Unit]
Description=Upstream DNS/HTTPS Probe
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/upstream-probe/secrets.env
ExecStart=/usr/bin/env upstream-probe --config /etc/upstream-probe/config.yml --write
Restart=on-failure
RestartSec=10
User=upstream-probe
Group=upstream-probe
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/upstream-probe

[Install]
WantedBy=multi-user.target
